# server-beta.service
[Unit]
Description=Service Beta

[Service]
WorkingDirectory=/home/
ExecStart=/home/server-beta

[Install]
WantedBy=multi-user.target