# server-beta.service
[Unit]
Description=Service Beta

[Service]
WorkingDirectory=/
ExecStart=/server-beta/bin/server-beta

[Install]
WantedBy=multi-user.target